# SolarProbe-Client/crontab
*/5 * * * *     /var/lib/slang/library/solarprobe/slang-client/client.sh
